# deploy/spotify-pi-thing.service
[Unit]
Description=Spotify Pi Thing API
After=network-online.target
WantedBy=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SpotifyPiThing
Environment=PATH=/home/pi/SpotifyPiThing/venv/bin
Environment=PYTHONPATH=/home/pi/SpotifyPiThing
ExecStart=/home/pi/SpotifyPiThing/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8888
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
